# Packs/SocFrameworkTrendMicroVisionOne/CorrelationRules/SOC Trend Micro Vision One V3.yml
rule_id: 0
fromversion: 6.10.0
tags:
  - SOCFramework
  - Passthrough
  - Endpoint
  - TrendMicroVisionOne
action: ALERTS
alert_category: User Defined
alert_description: $alert_description
alert_domain: DOMAIN_SECURITY
alert_fields:
  action_file_path: filepath
  action_file_sha256: sha256
  action_local_ip: local_ip
  action_remote_ip: remote_ip_str
  actor_effective_username: user_name
  actor_process_command_line: cmdline
  actor_process_image_name: filename
  actor_process_image_path: filepath
  actor_process_image_sha256: sha256
  agent_device_domain: domain
  agent_hostname: v1_host_name
  agent_id: v1_host_guid
  vendor: vendor_name
  product: product_name
  mitretacticid: mitre_tactic_id
  mitretacticname: mitre_tactic
  mitretechniqueid: mitre_ids_str
  mitretechniquename: mitre_ids_str
  alert_description: alert_description
  externallink: workbench_link
  external_pivot_url: workbench_link
  externalstatus: status
  originalalertid: id
  originalalertname: alert_name
  originalalertsource: alert_source
  severity: severity
  source_insert_ts: alert_time
  additionalindicators: ioc_value
  filehash: sha256
  mac: mac_address
  parentprocessname: parent_process_name
  parentprocesspath: parent_process_path
  processcmd: cmdline
  tim_main_indicator: ioc_value
  userid: user_id
  trendmicrovisiononexdrinvestigationstatus: investigation_status
  trendmicrovisiononexdrpriorityscore: score
alert_name: Trend Micro - $alert_name
alert_type: null
crontab: null
dataset: alerts
description: Creates an XSIAM alert for each Trend Micro Vision One workbench detection event
drilldown_query_timeframe: ALERT
execution_mode: REAL_TIME
global_rule_id: SOC Trend Micro Vision One V3
investigation_query_link: '// All activity from impacted host

  dataset = trend_micro_vision_one_v3_generic_alert_raw

  | filter agent_hostname = $v1_host_name

  | fields *

  '
is_enabled: true
lookup_mapping: []
mapping_strategy: CUSTOM
mitre_defs:
  TA0001 - Initial Access:
    - T1078 - Valid Accounts
    - T1190 - Exploit Public-Facing Application
    - T1566 - Phishing
  TA0002 - Execution:
    - T1047 - Windows Management Instrumentation
    - T1059 - Command and Scripting Interpreter
    - T1204 - User Execution
  TA0003 - Persistence:
    - T1053 - Scheduled Task/Job
    - T1136 - Create Account
    - T1547 - Boot or Logon Autostart Execution
    - T1574 - Hijack Execution Flow
  TA0004 - Privilege Escalation:
    - T1055 - Process Injection
    - T1134 - Access Token Manipulation
    - T1547 - Boot or Logon Autostart Execution
    - T1548 - Abuse Elevation Control Mechanism
  TA0005 - Defense Evasion:
    - T1027 - Obfuscated Files or Information
    - T1036 - Masquerading
    - T1055 - Process Injection
    - T1140 - Deobfuscate/Decode Files or Information
    - T1218 - System Binary Proxy Execution
    - T1562 - Impair Defenses
  TA0006 - Credential Access:
    - T1003 - OS Credential Dumping
    - T1110 - Brute Force
    - T1552 - Unsecured Credentials
    - T1558 - Steal or Forge Kerberos Tickets
  TA0007 - Discovery:
    - T1033 - System Owner/User Discovery
    - T1046 - Network Service Discovery
    - T1082 - System Information Discovery
    - T1083 - File and Directory Discovery
  TA0008 - Lateral Movement:
    - T1021 - Remote Services
    - T1550 - Use Alternate Authentication Material
  TA0009 - Collection:
    - T1005 - Data from Local System
    - T1039 - Data from Network Shared Drive
    - T1113 - Screen Capture
  TA0011 - Command and Control:
    - T1071 - Application Layer Protocol
    - T1095 - Non-Application Layer Protocol
    - T1105 - Ingress Tool Transfer
    - T1572 - Protocol Tunneling
  TA0010 - Exfiltration:
    - T1020 - Automated Exfiltration
    - T1041 - Exfiltration Over C2 Channel
    - T1567 - Exfiltration Over Web Service
  TA0040 - Impact:
    - T1485 - Data Destruction
    - T1486 - Data Encrypted for Impact
    - T1490 - Inhibit System Recovery
    - T1499 - Endpoint Denial of Service
name: SOC Trend Micro Vision One V3
search_window: null
severity: User Defined
simple_schedule: null
suppression_duration: 1 hours
suppression_enabled: true
suppression_fields:
  - id
timezone: null
user_defined_category: tactic
user_defined_severity: severity
xql_query: "dataset = trend_micro_vision_one_v3_generic_alert_raw\n| filter alert_provider\
  \ = \"SAE\"\n\n| alter j = _alert_data -> raw_json\n\n| alter j_str = to_string(j)\n\
  | alter mitre_technique_id_raw =\n    json_extract_scalar(j_str, \"$.matched_rules[0].matched_filters[0].mitre_technique_ids[0]\"\
  )\n| alter j_str = null\n\n| alter mitre_ids_str =\n    if(\n      mitre_technique_id_raw\
  \ != null and mitre_technique_id_raw != \"\",\n      if(mitre_technique_id_raw\
  \ contains \".\",\n         arrayindex(regextract(mitre_technique_id_raw, \"(T\\\
  \\d+)\\\\.\"), 0),\n         mitre_technique_id_raw),\n      \"-\"\n    )\n\n|\
  \ alter ta0043_reconnaissance       = arraycreate(\"T1590\",\"T1591\",\"T1592\"\
  ,\"T1593\",\"T1594\",\"T1595\",\"T1596\",\"T1597\",\"T1598\",\"T1599\")\n| alter\
  \ ta0042_resource_development = arraycreate(\"T1583\",\"T1584\",\"T1585\",\"T1586\"\
  ,\"T1587\",\"T1650\")\n| alter ta0001_initial_access       = arraycreate(\"T1078\"\
  ,\"T1189\",\"T1190\",\"T1195\",\"T1133\",\"T1200\",\"T1566\",\"T1091\")\n| alter\
  \ ta0002_execution            = arraycreate(\"T1059\",\"T1106\",\"T1047\",\"T1203\"\
  ,\"T1129\",\"T1559\",\"T1204\",\"T1072\")\n| alter ta0003_persistence          =\
  \ arraycreate(\"T1547\",\"T1543\",\"T1136\",\"T1505\",\"T1053\",\"T1078\",\"T1546\"\
  ,\"T1574\",\"T1037\")\n| alter ta0004_privilege_escalation = arraycreate(\"T1548\"\
  ,\"T1068\",\"T1078\",\"T1055\",\"T1134\",\"T1547\",\"T1574\")\n| alter ta0005_defense_evasion\
  \      = arraycreate(\"T1027\",\"T1070\",\"T1218\",\"T1140\",\"T1562\",\"T1036\"\
  ,\"T1055\",\"T1497\",\"T1620\",\"T1553\",\"T1222\",\"T1202\")\n| alter ta0006_credential_access\
  \    = arraycreate(\"T1003\",\"T1555\",\"T1552\",\"T1110\",\"T1621\",\"T1558\"\
  ,\"T1539\",\"T1606\")\n| alter ta0007_discovery            = arraycreate(\"T1082\"\
  ,\"T1083\",\"T1046\",\"T1057\",\"T1016\",\"T1049\",\"T1033\",\"T1518\",\"T1120\"\
  ,\"T1069\",\"T1087\",\"T1135\")\n| alter ta0008_lateral_movement     = arraycreate(\"\
  T1021\",\"T1210\",\"T1091\",\"T1072\",\"T1080\",\"T1550\")\n| alter ta0009_collection\
  \           = arraycreate(\"T1005\",\"T1039\",\"T1113\",\"T1114\",\"T1115\",\"\
  T1119\",\"T1074\",\"T1560\")\n| alter ta0011_command_and_control  = arraycreate(\"\
  T1071\",\"T1095\",\"T1105\",\"T1571\",\"T1572\",\"T1041\",\"T1132\",\"T1001\"\
  ,\"T1568\",\"T1573\")\n| alter ta0010_exfiltration         = arraycreate(\"T1041\"\
  ,\"T1567\",\"T1020\",\"T1048\",\"T1030\")\n| alter ta0040_impact               =\
  \ arraycreate(\"T1485\",\"T1486\",\"T1490\",\"T1499\",\"T1561\",\"T1529\",\"T1496\"\
  ,\"T1495\",\"T1491\")\n\n| alter mitre_tactic = if(ta0040_impact contains mitre_ids_str,\
  \ \"Impact\")\n| alter mitre_tactic = if(ta0010_exfiltration contains mitre_ids_str,\
  \ \"Exfiltration\", mitre_tactic)\n| alter mitre_tactic = if(ta0011_command_and_control\
  \ contains mitre_ids_str, \"Command and Control\", mitre_tactic)\n| alter mitre_tactic\
  \ = if(ta0009_collection contains mitre_ids_str, \"Collection\", mitre_tactic)\n\
  | alter mitre_tactic = if(ta0008_lateral_movement contains mitre_ids_str, \"Lateral\
  \ Movement\", mitre_tactic)\n| alter mitre_tactic = if(ta0007_discovery contains\
  \ mitre_ids_str, \"Discovery\", mitre_tactic)\n| alter mitre_tactic = if(ta0006_credential_access\
  \ contains mitre_ids_str, \"Credential Access\", mitre_tactic)\n| alter mitre_tactic\
  \ = if(ta0005_defense_evasion contains mitre_ids_str, \"Defense Evasion\", mitre_tactic)\n\
  | alter mitre_tactic = if(ta0004_privilege_escalation contains mitre_ids_str, \"\
  Privilege Escalation\", mitre_tactic)\n| alter mitre_tactic = if(ta0003_persistence\
  \ contains mitre_ids_str, \"Persistence\", mitre_tactic)\n| alter mitre_tactic\
  \ = if(ta0002_execution contains mitre_ids_str, \"Execution\", mitre_tactic)\n|\
  \ alter mitre_tactic = if(ta0001_initial_access contains mitre_ids_str, \"Initial\
  \ Access\", mitre_tactic)\n| alter mitre_tactic = if(ta0042_resource_development\
  \ contains mitre_ids_str, \"Resource Development\", mitre_tactic)\n| alter mitre_tactic\
  \ = if(ta0043_reconnaissance contains mitre_ids_str, \"Reconnaissance\", mitre_tactic)\n\
  \n| alter mitre_tactic_id = if(ta0040_impact contains mitre_ids_str, \"TA0040\"\
  )\n| alter mitre_tactic_id = if(ta0010_exfiltration contains mitre_ids_str, \"\
  TA0010\", mitre_tactic_id)\n| alter mitre_tactic_id = if(ta0011_command_and_control\
  \ contains mitre_ids_str, \"TA0011\", mitre_tactic_id)\n| alter mitre_tactic_id\
  \ = if(ta0009_collection contains mitre_ids_str, \"TA0009\", mitre_tactic_id)\n\
  | alter mitre_tactic_id = if(ta0008_lateral_movement contains mitre_ids_str, \"\
  TA0008\", mitre_tactic_id)\n| alter mitre_tactic_id = if(ta0007_discovery contains\
  \ mitre_ids_str, \"TA0007\", mitre_tactic_id)\n| alter mitre_tactic_id = if(ta0006_credential_access\
  \ contains mitre_ids_str, \"TA0006\", mitre_tactic_id)\n| alter mitre_tactic_id\
  \ = if(ta0005_defense_evasion contains mitre_ids_str, \"TA0005\", mitre_tactic_id)\n\
  | alter mitre_tactic_id = if(ta0004_privilege_escalation contains mitre_ids_str,\
  \ \"TA0004\", mitre_tactic_id)\n| alter mitre_tactic_id = if(ta0003_persistence\
  \ contains mitre_ids_str, \"TA0003\", mitre_tactic_id)\n| alter mitre_tactic_id\
  \ = if(ta0002_execution contains mitre_ids_str, \"TA0002\", mitre_tactic_id)\n|\
  \ alter mitre_tactic_id = if(ta0001_initial_access contains mitre_ids_str, \"TA0001\"\
  , mitre_tactic_id)\n| alter mitre_tactic_id = if(ta0042_resource_development contains\
  \ mitre_ids_str, \"TA0042\", mitre_tactic_id)\n| alter mitre_tactic_id = if(ta0043_reconnaissance\
  \ contains mitre_ids_str, \"TA0043\", mitre_tactic_id)\n\n| alter\n    id      \
  \             = j -> id,\n    status               = j -> status,\n    investigation_status\
  \ = j -> investigation_status,\n    investigation_result = j -> investigation_result,\n\
  \    workbench_link       = j -> workbench_link,\n    alert_provider       = j ->\
  \ alert_provider,\n    alert_name           = j -> model,\n    score           \
  \     = to_integer(j -> score),\n    severity             = j -> severity,\n  \
  \  alert_time           = j -> created_date_time,\n    alert_description    = j\
  \ -> description,\n    alert_source         = coalesce(j -> alert_provider, \"Trend\
  \ Micro Vision One\"),\n    indicators           = j -> indicators[]\n\n| alter\
  \ i_host = arrayindex(arrayfilter(indicators, json_extract_scalar(\"@element\",\"\
  $.type\") = \"host\"), 0)\n| alter\n    v1_host_guid = json_extract_scalar(i_host,\
  \ \"$.value.guid\"),\n    v1_host_name = json_extract_scalar(i_host, \"$.value.name\"\
  ),\n    local_ip     = replace(json_extract_scalar(i_host, \"$.value.ips[0]\"),\
  \ \"\\\\\"\", \"\")\n\n| alter mac_address =\n    coalesce(\n      json_extract_scalar(i_host,\
  \ \"$.value.mac\"),\n      json_extract_scalar(i_host, \"$.value.mac_address\"\
  ),\n      json_extract_scalar(i_host, \"$.value.macs[0]\"),\n      json_extract_scalar(i_host,\
  \ \"$.value.macAddresses[0]\")\n    )\n\n| alter i_user = arrayindex(arrayfilter(indicators,\
  \ json_extract_scalar(\"@element\",\"$.type\") = \"user_account\"), 0)\n| alter\n\
  \    user_name = json_extract_scalar(i_user, \"$.value\"),\n    user_id   = null\n\
  \n| alter i_cmd1 = arrayindex(arrayfilter(indicators, json_extract_scalar(\"@element\"\
  ,\"$.type\") = \"command_line\"), 0)\n| alter i_cmd2 = arrayindex(arrayfilter(indicators,\
  \ json_extract_scalar(\"@element\",\"$.field\") = \"processCmd\"), 0)\n| alter\
  \ cmdline = coalesce(json_extract_scalar(i_cmd1, \"$.value\"), json_extract_scalar(i_cmd2,\
  \ \"$.value\"))\n\n| alter i_sha = arrayindex(arrayfilter(indicators, json_extract_scalar(\"\
  @element\",\"$.type\") = \"file_sha256\"), 0)\n| alter sha256 = json_extract_scalar(i_sha,\
  \ \"$.value\")\n\n| alter i_peer = arrayindex(arrayfilter(indicators, json_extract_scalar(\"\
  @element\",\"$.field\") = \"peerIp\"), 0)\n| alter remote_ip_str = json_extract_scalar(i_peer,\
  \ \"$.value\")\n\n| alter i_dom = arrayindex(arrayfilter(indicators, json_extract_scalar(\"\
  @element\",\"$.field\") = \"domain\"), 0)\n| alter domain = json_extract_scalar(i_dom,\
  \ \"$.value\")\n\n| alter i_pfp = arrayindex(arrayfilter(indicators, json_extract_scalar(\"\
  @element\",\"$.field\") = \"parentFilePath\"), 0)\n| alter parent_process_path\
  \ = json_extract_scalar(i_pfp, \"$.value\")\n| alter parent_process_name = replace(parent_process_path,\
  \ \"^.*[\\\\\\\\/]\", \"\")\n\n| alter i_reg = arrayindex(arrayfilter(indicators,\
  \ json_extract_scalar(\"@element\",\"$.field\") = \"objectRegistryData\"), 0)\n\
  | alter reg_path = json_extract_scalar(i_reg, \"$.value\")\n| alter filepath =\n\
  \    coalesce(\n      reg_path,\n      arrayindex(regextract(cmdline, \"^\\\\s*([^\\\
  \\s]+)\"), 0)\n    )\n| alter filename = replace(filepath, \"^.*[\\\\\\\\/]\", \"\
  \")\n\n| alter ioc_value = coalesce(sha256, null)\n\n| alter vendor_name = \"Trend Micro\", product_name = \"Vision One\"\n\n| fields\n    id, workbench_link,\
  \ alert_name, alert_source, status,\n    investigation_status, investigation_result,\n\
  \    score, severity, alert_time, alert_description,\n    v1_host_guid, v1_host_name,\
  \ local_ip, mac_address,\n    user_name, user_id,\n    filename, filepath, parent_process_path,\
  \ parent_process_name, cmdline,\n    sha256, ioc_value, domain, remote_ip_str,\n\
  \    mitre_ids_str, mitre_tactic, mitre_tactic_id\n"
